# .package/go-tcp-server/DEBIAN/postinst
#!/usr/bin/env bash

set -e

if [[ ! -f "/etc/go-tcp-server/config.yml" ]]; then
  cp -a "/etc/go-tcp-server/config.example.yml" "/etc/go-tcp-server/config.yml"
fi

systemctl daemon-reload && \
    systemctl start go-tcp-server && \
    systemctl enable go-tcp-server
