# .package/go-tcp-server/etc/systemd/system/go-tcp-server.service
[Unit]
Description=TCP server written in Go
After=network.target
After=systemd-user-sessions.service
After=network-online.target

[Service]
Type=simple
User=go-tcp-server
ExecStart=/usr/local/bin/go-tcp-server --config /etc/go-tcp-server/config.yml
TimeoutSec=10
Restart=on-failure
RestartSec=30
StartLimitInterval=350
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
